# docker/periodic/hourly/zipkin-dependencies-today
#!/bin/sh
java ${JAVA_OPTS} -jar /zipkin-dependencies/zipkin-dependencies.jar
