# docker/periodic/hourly/zipkin-dependencies-today
#!/bin/sh
#
# Copyright The OpenZipkin Authors
# SPDX-License-Identifier: Apache-2.0
#

cd /zipkin-dependencies
start-zipkin-dependencies
